# health_check.conf
description "health_check"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel [016]

console log

exec 
